# templates/ai-runtime@.service
[Unit]
Description=AI Runtime for %i
After=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
Environment=HOME=/home/%i
Environment=AI_RUNTIME_USER=%i
ExecStart=/home/%i/runtime/start.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Allow sudo for privileged runtimes
NoNewPrivileges=no
ProtectSystem=strict
ProtectHome=off
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
